# ops/deploy/adarkwa-bot-deploy.service
[Unit]
Description=Deploy the latest Adarkwa Study Bot production image
After=network-online.target k3s.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/adarkwa-study-bot/deploy.env
ExecStart=/usr/local/bin/adarkwa-bot-deploy.sh
